# docker/app/supervisor/conf.d/background.conf
[program:wp_cron_event]
command=/bin/sh -c "/usr/local/bin/wp cron event run --due-now && sleep 60"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3

#  [program:wp_worker]
#  command=/usr/local/bin/wp background-worker listen 
#  directory=/var/www/html
#  stdout_logfile=/dev/stdout
#  stdout_logfile_maxbytes=0
#  stderr_logfile=/dev/stderr
#  stderr_logfile_maxbytes=0
#  autorestart=false